# vireon-release/vps-control-plane/systemd/vireon-indexer-refresh.service
[Unit]
Description=Refresh Veiron local index
After=veiron-node.service

[Service]
Type=oneshot
User=veiron
Group=veiron
WorkingDirectory=/opt/veiron
# `sync` only rebuilds when chain tip diverged (cheaper + reorg-safe).
ExecStart=/opt/veiron/bin/veiron-indexer --network mainnet-candidate --chain-data-dir /var/lib/veiron/.veiron-mainnet/chain --index-dir /var/lib/veiron/.veiron-mainnet/indexer sync
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=/var/lib/veiron/.veiron-mainnet/chain
ReadWritePaths=/var/lib/veiron/.veiron-mainnet/indexer
MemoryMax=1G
CPUQuota=50%
CPUWeight=25
IOWeight=25
Nice=10
