# vireon-release/vps-control-plane/systemd/vireon-rpc.service
[Unit]
Description=Veiron Mainnet Candidate RPC gateway
After=veiron-node.service
Requires=veiron-node.service

[Service]
Type=simple
User=veiron
Group=veiron
WorkingDirectory=/opt/veiron
ExecStart=/opt/veiron/bin/veiron-rpc-gateway --config /etc/veiron/rpc.toml --node-config /etc/veiron/node.toml
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
# Chain must be writable: public-submit mining posts append blocks via RPC.
ReadOnlyPaths=/var/lib/veiron/.veiron-mainnet/indexer /var/lib/veiron/.veiron-mainnet/node
ReadWritePaths=/var/lib/veiron/.veiron-mainnet/mempool /var/lib/veiron/.veiron-mainnet/chain
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
# Keep enough headroom for a cache refresh while applying memory pressure before OOM.
MemoryMax=3G
MemoryHigh=2G
OOMPolicy=continue
CPUQuota=100%

[Install]
WantedBy=multi-user.target
