# vireon-release/vps-control-plane/systemd/vireon-vps-admin.service
[Unit]
Description=Veiron VPS fleet administration agent
After=network-online.target veiron-node.service veiron-rpc.service
Wants=network-online.target

[Service]
Type=simple
User=veiron-control
Group=veiron-control
WorkingDirectory=/opt/veiron
ExecStart=/opt/veiron/bin/veiron-vps-admin --config /etc/veiron-control/admin.toml
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/veiron-control
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
LockPersonality=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target

