# vireon-release/vps-control-plane/systemd/vireon-auto-update.service
[Unit]
Description=Veiron auto-update from GitHub Releases (node/rpc/indexer/admin/pool)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
ExecStart=/opt/veiron/vps-control-plane/auto-update.sh
Nice=10
TimeoutStartSec=900
# No interactive approval — apply the newest matching GitHub release immediately.
Environment=VEIRON_GITHUB_REPO=andreidohot/veiron-network

[Install]
WantedBy=multi-user.target
