# vireon-release/vps-control-plane/systemd/vireon-mining-pool.service
[Unit]
Description=Veiron pooled mining coordinator prototype
After=network-online.target veiron-node.service veiron-rpc.service
Requires=veiron-node.service veiron-rpc.service

[Service]
Type=simple
User=veiron-pool
Group=veiron-pool
WorkingDirectory=/opt/veiron
ExecStart=/opt/veiron/bin/veiron-mining-pool --config /etc/veiron-pool/pool.toml
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/veiron-pool
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryMax=1G
CPUQuota=200%

[Install]
WantedBy=multi-user.target

