# vireon-release/vps-control-plane/systemd/vireon-node.service
[Unit]
Description=Veiron Mainnet Candidate full node
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=veiron
Group=veiron
WorkingDirectory=/opt/veiron
ExecStart=/opt/veiron/bin/veiron-node --config /etc/veiron/node.toml --data-dir /var/lib/veiron/.veiron-mainnet/chain --mempool-dir /var/lib/veiron/.veiron-mainnet/mempool start-node
ExecStop=/opt/veiron/bin/veiron-node --config /etc/veiron/node.toml --data-dir /var/lib/veiron/.veiron-mainnet/chain --mempool-dir /var/lib/veiron/.veiron-mainnet/mempool shutdown
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/veiron
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryMax=3G
MemoryHigh=2G
OOMPolicy=continue
CPUQuota=200%
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
