# veiron-release/vps-control-plane/systemd/veiron-auto-update.timer
[Unit]
Description=Poll GitHub for Veiron VPS control-plane updates
After=network-online.target

[Timer]
# First probe shortly after boot, then every 15 minutes.
OnBootSec=2min
OnUnitActiveSec=15min
AccuracySec=1min
Persistent=true
Unit=veiron-auto-update.service

[Install]
WantedBy=timers.target
