# veiron-release/vps-control-plane/systemd/veiron-indexer-refresh.timer
[Unit]
Description=Refresh Veiron local index every 30 seconds

[Timer]
OnBootSec=30
OnUnitActiveSec=30
AccuracySec=5
RandomizedDelaySec=3
Persistent=true

[Install]
WantedBy=timers.target
